# contrib/deb/postinst
#!/bin/sh
set -e

# Create system group if sysusers.d didn't handle it
if ! getent group yggdrasil >/dev/null 2>&1; then
    groupadd --system yggdrasil
fi

# Create config directory
mkdir -p /etc/yggdrasil
chmod 0750 /etc/yggdrasil

# Generate default config if not present, otherwise normalize the existing
# one so any new fields introduced by this version get spliced in (with
# their explanatory comments) while user values and comments are preserved.
if [ ! -f /etc/yggdrasil/yggdrasil.toml ]; then
    /usr/bin/yggdrasil --genconf=/etc/yggdrasil/yggdrasil.toml
    chmod 0640 /etc/yggdrasil/yggdrasil.toml
    chgrp yggdrasil /etc/yggdrasil/yggdrasil.toml
    echo ""
    echo "A fresh config has been generated at /etc/yggdrasil/yggdrasil.toml"
    echo "Edit it to add peers (https://publicpeers.neilalexander.dev), then enable and start the service:"
    echo "  sudo systemctl enable --now yggdrasil"
    echo ""
else
    BACKUP=/var/backups/yggdrasil.toml.$(date +%Y%m%d%H%M%S)
    mkdir -p /var/backups
    cp -p /etc/yggdrasil/yggdrasil.toml "$BACKUP"
    if /usr/bin/yggdrasil --normalize="$BACKUP" > /etc/yggdrasil/yggdrasil.toml.new; then
        chmod 0640 /etc/yggdrasil/yggdrasil.toml.new
        chgrp yggdrasil /etc/yggdrasil/yggdrasil.toml.new
        mv /etc/yggdrasil/yggdrasil.toml.new /etc/yggdrasil/yggdrasil.toml
    else
        echo "WARNING: failed to normalize /etc/yggdrasil/yggdrasil.toml (kept as-is, backup at $BACKUP)" >&2
        rm -f /etc/yggdrasil/yggdrasil.toml.new
    fi
fi

# Reload systemd
systemctl daemon-reload || true

#DEBHELPER#
